# result_controller.service
[Unit]
Description=Allure result controller
After=multi-user.target
Conflicts=getty@tty1.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/result-controller/src/ResultController.py
StandardInput=tty-force

[Install]
WantedBy=multi-user.target
